# sample-app/.crontab
*/10 * * * * * pwd
*/10 * * * * * ls -al